# deploy/booty.service
[Unit]
Description=Booty - GitHub issue to PR agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=booty
Group=www-data
WorkingDirectory=${INSTALL_DIR}
EnvironmentFile=-${INSTALL_DIR}/.env
EnvironmentFile=/etc/booty/release.env
EnvironmentFile=-/etc/booty/secrets.env
ExecStart=${INSTALL_DIR}/.venv/bin/uvicorn booty.main:app --uds ${INSTALL_DIR}/booty.sock
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
